# root/etc/update-motd.d/80-esm
#!/bin/sh

. /etc/RGBux/system.conf

cat <<EOF

System configured to use a 15Khz CRT on ${CONNECTOR_DRM}. Feel the feeling!
EOF
echo
