# root/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
#ExecStart=-/sbin/agetty --noissue --autologin arcade %I $TERM
ExecStart=-/sbin/agetty --noissue -n -l /opt/RGBux/bin/system-init %I $TERM
Type=idle