# root/etc/rc.local
#!/bin/bash

# RGBux scripts on init
run-parts /opt/RGBux/autostart

exit 0