# root/etc/rc.local
#!/bin/bash

# RGBux scripts on init
run-parts /opt/RGBux/autostart

# Force dhcp ip
dhclient

exit 0
